# dotfiles/.config/autostart/fancurve.desktop
[Desktop Entry]
Name=nan0s7's fan curve script
Type=Application
Path=/home/dalisoft/nfancurve
Exec=/home/dalisoft/nfancurve/update.sh